# apps/gotosocial/root/etc/s6-overlay/s6-rc.d/gotosocial-cli/run
#!/command/with-contenv ash

# Launch appnest-terminal
/gotosocial-cli.sh